# deployment/systemd/sentinel.service
# Systemd service file for SENTINEL
# Location: /etc/systemd/system/sentinel.service
#
# Installation:
# 1. sudo cp deployment/systemd/sentinel.service /etc/systemd/system/
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl enable sentinel
# 4. sudo systemctl start sentinel

[Unit]
Description=SENTINEL Apartment Controller API
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=sentinel
Group=sentinel
WorkingDirectory=/opt/sentinel

# Environment
Environment="PATH=/opt/sentinel/venv/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/sentinel/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sentinel/data

# Process management
ExecStart=/opt/sentinel/venv/bin/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --log-config /opt/sentinel/logging.conf

ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel

[Install]
WantedBy=multi-user.target
